# debian/nebula-server.service
[Unit]
Description=nebula-server

[Service]
ExecStart=/usr/bin/nebula-server

[Install]
WantedBy=first-app.service
