# debian/nebula-metrics.service
[Unit]
Description=nebula-monitor

[Service]
ExecStart=/usr/bin/nebula-monitor

[Install]
WantedBy=multi-user.target
